# deploy/vrops-api.service
[Unit]
Description=VR-OPS RAG FastAPI
After=network.target

[Service]
Type=simple
User=vrops
WorkingDirectory=/opt/vr-ops-rag
EnvironmentFile=/opt/vr-ops-rag/.env
ExecStart=/opt/vr-ops-rag/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
