# deploy/vrops-dashboard.service
[Unit]
Description=VR-OPS RAG Streamlit Dashboard
After=network.target vrops-api.service

[Service]
Type=simple
User=vrops
WorkingDirectory=/opt/vr-ops-rag
EnvironmentFile=/opt/vr-ops-rag/.env
ExecStart=/opt/vr-ops-rag/.venv/bin/streamlit run dashboard/app.py \
    --server.port 8051 \
    --server.address 0.0.0.0 \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
